# roles/service_spamassassin/files/spamassassin-update.timer
[Unit]
Description=Spamassassin Updates

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target